# 05_VideoAnalysis/connector/systemd/connector.service
[Unit]
Description=HomeNet RTSP connector service
After=docker.service

[Service]
ExecStart=/app/run-connector.sh

[Install]
WantedBy=default.target